# blocky-dev/rootfs/etc/services.d/blocky/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Handle Blocky service shutdown
# Take down the S6 supervision tree when blocky fails
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

declare APP_EXIT_CODE=${1}

# Exit code handling:
# - 0: Clean shutdown (normal), allow restart
# - >256: Signal-derived (256 + signal number, e.g., SIGTERM = 271), allow restart
# - 1-256: Application error or abnormal exit (256 = abnormal without signal), halt

if [[ "${APP_EXIT_CODE}" -eq 0 ]]; then
  bashio::log.info "Service stopped cleanly (exit code 0)"
  exit 0
fi

if [[ "${APP_EXIT_CODE}" -gt 256 ]]; then
  bashio::log.info "Service stopped via signal-derived exit code ${APP_EXIT_CODE}"
  exit 0
fi

bashio::log.warning "Halt add-on with exit code ${APP_EXIT_CODE}"
bashio::log.warning "Blocky exited unexpectedly. Check logs for errors."

echo "${APP_EXIT_CODE}" > /run/s6-linux-init-container-results/exitcode
exec /run/s6/basedir/bin/halt
