# blocky/rootfs/etc/services.d/blocky/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start blocky as service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================
CONFIG="/etc/blocky.yaml"

exec blocky --config "${CONFIG}"
